# ansible/templates/fleetmesh.service.j2
[Unit]
Description=FleetMesh ship agent ({{ fleetmesh_ship_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ fleetmesh_service_user }}
WorkingDirectory={{ fleetmesh_app_dir }}
ExecStart={{ fleetmesh_node_bin }} {{ fleetmesh_app_dir }}/bin/fleetmesh.js --config {{ fleetmesh_config_dir }}/ship.config.json --creds {{ fleetmesh_config_dir }}/.tgcreds.json
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
